# supervisor/command.ini
[program:statsd-notify]
command=/usr/bin/php vendor/bin/statsd-console statsd:notify /app/gauges.yml --statsd-host localhost --statsd-port 8125 --verbose
directory=/app
startsecs=1

user=root

redirect_stderr=false

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
